# supervisord-icecast.ini
[program:icecast]
command=/usr/bin/icecast -c /etc/icecast.xml
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
